# American Betting League/com.abl.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.abl.monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rico/Library/Mobile Documents/com~apple~CloudDocs/Really Rico /Python Coding/American Betting League/run_monitor.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rico/Library/Mobile Documents/com~apple~CloudDocs/Really Rico /Python Coding/American Betting League/monitor_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rico/Library/Mobile Documents/com~apple~CloudDocs/Really Rico /Python Coding/American Betting League/monitor_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/rico/Library/Mobile Documents/com~apple~CloudDocs/Really Rico /Python Coding/American Betting League</string>
</dict>
</plist>
